# overlay/etc/init.d/displayip
#!/sbin/openrc-run

command="python3 ~/.setup/displayip.py"
